# deploy/systemd/parser-news-docker.service.example
# ParserNews: автозапуск Docker Compose при загрузке сервера (без привязки к SSH).
#
# Установка:
#   1) Подставьте свой путь к проекту в WorkingDirectory (ниже — /opt/parser-news).
#   2) Проверьте, что команда docker доступна: which docker compose  или  which docker-compose
#   3)
#      sudo cp deploy/systemd/parser-news-docker.service.example /etc/systemd/system/parser-news.service
#      sudo nano /etc/systemd/system/parser-news.service
#      sudo systemctl daemon-reload
#      sudo systemctl enable parser-news.service
#      sudo systemctl start parser-news.service
#
# Статус и логи:
#   sudo systemctl status parser-news
#   journalctl -u parser-news -f

[Unit]
Description=ParserNews (Docker Compose)
# Сначала Docker и сеть, иначе compose up может не найти сокет / сеть
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/parser-news
# Поднимает контейнеры в фоне; после выхода SSH они продолжают работать
ExecStart=/usr/bin/docker compose up -d
# При остановке юнита (systemctl stop, выключение) — мягко останавливаем стек
# (тома сохраняются, если не использовали -v)
ExecStop=/usr/bin/docker compose stop
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
